# docker/config/supervisor/supervisord.conf
; =============================================================================
; AgentOS Daemon Supervisor 配置
; 版本: 1.0.0
; 用途: 管理所有守护进程的生命周期
;
; 守护进程列表:
;   - llm_d: LLM服务管理守护进程
;   - market_d: 工具市场服务守护进程
;   - monit_d: 系统监控守护进程
;   - sched_d: 任务调度守护进程
;   - tool_d: 工具执行守护进程
;   - gateway_d: API网关守护进程
; =============================================================================

[supervisord]
; 使用前台模式运行（Docker容器需要）
nodaemon=true

; 日志配置
logfile=/app/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info

; PID文件
pidfile=/app/run/supervisord.pid

; 子进程目录
directory=/app

; 用户和组
user=agentos

; UNIX socket权限（用于supervisorctl）
chmod=0770
chown=agentos:agentos

; HTTP服务器（可选，用于Web界面）
;[inet_http_server]
;port=*:9001
;username=admin
;password=your_password_here

; =============================================================================
; 各守护进程配置
; =============================================================================

[program:llm_d]
command=python3 -m agentos.daemon.llm_d --config /app/config/llm.yaml
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
user=agentos
redirect_stderr=true
stdout_logfile=/app/logs/llm_d.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
environment=AGENTOS_HOME="/app",LOG_LEVEL="%(ENV_LOG_LEVEL)s"

[program:market_d]
command=python3 -m agentos.daemon.market_d --config /app/config/market.yaml
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
user=agentos
redirect_stderr=true
stdout_logfile=/app/logs/market_d.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
environment=AGENTOS_HOME="/app",LOG_LEVEL="%(ENV_LOG_LEVEL)s"

[program:monit_d]
command=python3 -m agentos.daemon.monit_d --config /app/config/monit.yaml
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
user=agentos
redirect_stderr=true
stdout_logfile=/app/logs/monit_d.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
environment=AGENTOS_HOME="/app",LOG_LEVEL="%(ENV_LOG_LEVEL)s"

[program:sched_d]
command=python3 -m agentos.daemon.sched_d --config /app/config/sched.yaml
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
user=agentos
redirect_stderr=true
stdout_logfile=/app/logs/sched_d.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
environment=AGENTOS_HOME="/app",LOG_LEVEL="%(ENV_LOG_LEVEL)s"

[program:tool_d]
command=python3 -m agentos.daemon.tool_d --config /app/config/tool.yaml
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
user=agentos
redirect_stderr=true
stdout_logfile=/app/logs/tool_d.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
environment=AGENTOS_HOME="/app",LOG_LEVEL="%(ENV_LOG_LEVEL)s"

[program:gateway_d]
command=python3 -m agentos.daemon.gateway_d --config /app/config/gateway.yaml
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
user=agentos
redirect_stderr=true
stdout_logfile=/app/logs/gateway_d.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
environment=AGENTOS_HOME="/app",LOG_LEVEL="%(ENV_LOG_LEVEL)s"

; =============================================================================
; 进程组定义（便于批量管理）
; =============================================================================

[group:daemon-services]
programs=llm_d,market_d,monit_d,sched_d,tool_d,gateway_d
priority=999
